# debian/b3-utils.bubba-buttond.service
[Unit]
Description=B3 buttond service

[Service]
Type=forking
PIDFile=/run/bubba-buttond.pid
ExecStart=/usr/sbin/bubba-buttond

[Install]
WantedBy=multi-user.target
